# alpine/crontab
## Every minute print the date
* * * * * date
